# dotconfig/autostart/kDrive.desktop
[Desktop Entry]
Name=kDrive
GenericName=File Synchronizer
Exec="/usr/bin/kDrive"
Terminal=false
Icon=kdrive
Categories=Network
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=10
